# cicdProject.service
[Unit]
Description=Learn Cicd

[Service]
Type=simple
Restart=always
ExecStart= /home/adam/cicd

[Install]
WantedBy=multi-user.target
